# lab40-advanced-troubleshooting-recovery/scripts/problem-service.service
[Unit]
Description=Problem Service for Testing
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c "while true; do echo 'Service running'; sleep 30; done"
ExecStop=/bin/bash -c "exit 1"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
